# ree-runner.service
# REE Experiment Runner -- systemd service template for cloud workers
#
# Installation (on the cloud VM):
#   1. Copy this file, substituting BOTH placeholders:
#        sed "s/<user>/ree/g; s/<machine>/ree-cloud-N/g" ree-runner.service \
#          | sudo tee /etc/systemd/system/ree-runner.service > /dev/null
#      (replace ree-cloud-N with the machine name matching
#       REE_assembly/contributors/machines/ and the queue's machine_affinity)
#   2. Enable and start:
#        sudo systemctl daemon-reload
#        sudo systemctl enable ree-runner
#        sudo systemctl start ree-runner
#
# The runner starts automatically on boot (for auto-scaler integration).
# On shutdown, systemd sends SIGTERM which triggers the runner's graceful drain
# (finishes current experiment, releases claim, then exits).
#
# Logs: journalctl -u ree-runner -f

[Unit]
Description=REE Experiment Runner (auto-sync cloud worker)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=<user>
WorkingDirectory=/home/<user>/Documents/GitHub/REE_Working/ree-v3
ExecStart=/home/<user>/.venv/ree/bin/python3 experiment_runner.py --auto-sync --loop --loop-interval 120 --machine <machine>

# Graceful drain: give 10 minutes for current experiment episode to finish
TimeoutStopSec=600
KillSignal=SIGTERM

# Do not auto-restart -- the GitHub Actions auto-scaler manages lifecycle
Restart=no

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ree-runner

[Install]
WantedBy=multi-user.target
