# ree-runner.service
# REE Experiment Runner -- systemd service template for cloud workers
#
# Installation (on the cloud VM):
#   1. Copy this file, substituting BOTH placeholders:
#        sed "s/<user>/ree/g; s/<machine>/ree-cloud-N/g" ree-runner.service \
#          | sudo tee /etc/systemd/system/ree-runner.service > /dev/null
#      (replace ree-cloud-N with the machine name matching
#       REE_assembly/contributors/machines/ and the queue's machine_affinity)
#   2. Enable and start:
#        sudo systemctl daemon-reload
#        sudo systemctl enable ree-runner
#        sudo systemctl start ree-runner
#
# The runner starts automatically on boot (for auto-scaler integration).
# On shutdown, systemd sends SIGTERM which triggers the runner's graceful drain
# (finishes current experiment, releases claim, then exits).
#
# Logs: journalctl -u ree-runner -f

[Unit]
Description=REE Experiment Runner (auto-sync cloud worker)
After=network-online.target
Wants=network-online.target
# StartLimit* belong in [Unit] in modern systemd (>= v229). See [Service]
# Restart= block below for the rationale on these specific values.
StartLimitBurst=5
StartLimitIntervalSec=900

[Service]
Type=simple
User=<user>
WorkingDirectory=/home/<user>/REE_Working/ree-v3
ExecStart=/home/<user>/.venv/ree/bin/python3 experiment_runner.py --auto-sync --loop --loop-interval 120 --remote-control --machine <machine>

# Graceful drain: give 10 minutes for current experiment episode to finish
TimeoutStopSec=600
KillSignal=SIGTERM

# Auto-restart on non-zero exit, capped to prevent infinite tight loop.
# Rationale: the 2026-05-27 ree-cloud-1/2/3 fleet wedge -- V3-EXQ-610 queue entry
# was pushed without its script (producer-side atomicity violation in commit
# a7888df), validate_queue.py crashed all three runners with exit code 1, and
# the prior Restart=no policy left them dead for ~24h until manual restart.
# Restart=on-failure means a transient bad-queue commit self-heals in ~60s
# when the producer pushes the missing script. The StartLimit* cap means a
# persistent bad commit fails the unit loudly within 15min (5 retries x 60s),
# so heartbeat freshness alarms still fire instead of getting masked by
# infinite restart churn.
# SIGTERM (graceful drain on shutdown) is NOT treated as failure; the
# auto-scaler still owns clean lifecycle transitions.
Restart=on-failure
RestartSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ree-runner

[Install]
WantedBy=multi-user.target
